# borg/.config/systemd/user/borgmatic.service
[Unit]
Description=backup home directory with borgmatic


[Service]
Type=oneshot
ExecStart=/bin/sh -c 'borgmatic'
