# borg/.config/systemd/user/borgmatic.timer
[Unit]
Description=Run borgmatic hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
