# borg/.config/systemd/user/borg-backup.timer
[Unit]
Description=Run borg-backup daily

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
